# scripts/com.praxis.agent-sessions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.praxis.agent-sessions</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/nate/Praxis/Code&amp;DBs/Workflow/surfaces/api/agent_sessions.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/nate/Praxis</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/nate/Praxis/Code&amp;DBs/Workflow</string>
        <key>PATH</key>
        <string>/Users/nate/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/nate/Praxis/artifacts/agent_sessions.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/nate/Praxis/artifacts/agent_sessions.err</string>
</dict>
</plist>
