# scripts/com.praxis.agent-sessions.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.praxis.agent-sessions</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PRAXIS_PYTHON_BIN__</string>
        <string>__PRAXIS_REPO_ROOT__/Code&amp;DBs/Workflow/surfaces/api/agent_sessions.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PRAXIS_REPO_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__PRAXIS_REPO_ROOT__/Code&amp;DBs/Workflow</string>
        <key>PATH</key>
        <string>__PRAXIS_PATH__</string>
    </dict>

    <key>KeepAlive</key>
    <true/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/agent_sessions.log</string>

    <key>StandardErrorPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/agent_sessions.err</string>
</dict>
</plist>
