# scripts/com.praxis.public-branch-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.praxis.public-branch-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PRAXIS_REPO_ROOT__/scripts/public_branch_sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PRAXIS_REPO_ROOT__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>4</integer>
        <key>Minute</key>
        <integer>45</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/public_branch_sync.log</string>

    <key>StandardErrorPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/public_branch_sync.err</string>
</dict>
</plist>
