# scripts/com.praxis.authority-memory-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.praxis.authority-memory-refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>__PRAXIS_PYTHON_BIN__</string>
        <string>__PRAXIS_REPO_ROOT__/scripts/refresh_authority_memory.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PRAXIS_REPO_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>__PRAXIS_REPO_ROOT__/Code&amp;DBs/Workflow</string>
        <key>WORKFLOW_DATABASE_URL</key>
        <string>__WORKFLOW_DATABASE_URL__</string>
    </dict>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/authority_memory_refresh.log</string>

    <key>StandardErrorPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/authority_memory_refresh.err</string>
</dict>
</plist>
