# scripts/com.praxis.agent-session-spawn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.praxis.agent-session-spawn</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__PRAXIS_REPO_ROOT__/scripts/run_agent_session_workflow.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__PRAXIS_REPO_ROOT__</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>WORKFLOW_DATABASE_URL</key>
        <string>__WORKFLOW_DATABASE_URL__</string>
        <key>PATH</key>
        <string>__PRAXIS_PATH__</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/agent_session_spawn.log</string>

    <key>StandardErrorPath</key>
    <string>__PRAXIS_REPO_ROOT__/artifacts/agent_session_spawn.err</string>
</dict>
</plist>
